# confs/bootup.timer
[Unit]
Description=Run this at boot

[Timer]
OnBootSec=1min

[Install]
WantedBy=timers.target
