# confs/rc.local
#!/bin/bash
cp /usr/bin/resolv.conf /etc/
ip link set wlp2s0 down
systemctl stop postfix portsentry
exit 0
